# build/supervisor/conf.d/app.conf
[program:basil_worker_handler_consume_messages]
command = composer test-integration-listen
stopsignal=QUIT
